# systemd/transip-ddns.timer
[Unit]
Description=Run TransIP DDNS every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
